# deployment/scheduler.cron
# ⚠️ مسار نشر يدويّ (bare-metal) فقط — غير مستخدَم في نشر Docker/Coolify النشط
#    (الجدولة هناك عبر خدمة scheduler في docker-compose.yml). يُحفَظ مرجعًا لنشر غير-Docker.
#    انظر deployment/DEPLOYMENT.md.
#
# جدولة Laravel — سطر واحد فقط (يوزّع كل المهام + health:check)
# ثبّته عبر: crontab -u www-data -e
# ⚠ بدون هذا السطر: لا نسخ احتياطي، لا تنظيف، لا فحوصات صحّة.

* * * * * cd /var/www/alphacms && /usr/bin/php artisan schedule:run >> /dev/null 2>&1
