# deployment/supervisor-queue-worker.conf
; ⚠️ مسار نشر يدويّ (bare-metal) فقط — غير مستخدَم في نشر Docker/Coolify النشط
;    (الطابور هناك عبر خدمتَي worker/worker-media في docker-compose.yml، مسار /var/www/html).
;    يُحفَظ مرجعًا لنشر غير-Docker. انظر deployment/DEPLOYMENT.md.
;
; Supervisor (بديل systemd) — عامل طابور AlphaCMS
; انسخ إلى /etc/supervisor/conf.d/alphacms-queue.conf ثم:
;   supervisorctl reread && supervisorctl update && supervisorctl start alphacms-queue:*

[program:alphacms-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/alphacms/artisan queue:work redis --queue=default,cdn-purge --tries=3 --backoff=10 --max-time=3600 --sleep=3 --timeout=60
autostart=true
autorestart=true
stopwaitsecs=3600
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/alphacms/storage/logs/queue-worker.log
stopasgroup=true
killasgroup=true
